# packages/awscli/v2.29.0/profile.d/01-set-environ.sh
#!/bin/bash

# Setup AWS CLI autocompletion
if ! command -v aws >/dev/null 2>&1; then
    return
fi

complete -C '/usr/local/bin/aws_completer' aws